# ti_wisun_webapp/AM64x_wisun_webapp_service/95-serialdaemon-remove.rules
ACTION=="remove", KERNEL=="ttyACM0", SUBSYSTEM=="tty", TAG+="systemd", RUN+="/bin/systemctl --no-block start ti_wisunfan_app_stop.service"
